# gpu-power.service
[Unit]
Description=Power down the dGPU in favor of the iGPU

[Service]
ExecStart=/usr/bin/gpu-power

[Install]
WantedBy=multi-user.target
